# tailscale-receive-go.service
[Unit]
Description=Tailscale File Receiver Service (Go)
Documentation=https://github.com/1999AZZAR/tailscale_receiver
After=network-online.target tailscale.service
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
EnvironmentFile=-/etc/default/tailscale-receive
ExecStart=/usr/local/bin/tailscale-receiver-go
Restart=always
RestartSec=10s

# Device & file system isolation
PrivateTmp=true
ProtectSystem=strict
ProtectHome=no
ProtectClock=true
ProtectHostname=true
ProtectControlGroups=true
ProtectKernelTunables=true

# Privilege sandboxing
NoNewPrivileges=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictRealtime=true
CapabilityBoundingSet=CAP_CHOWN CAP_DAC_OVERRIDE CAP_NET_RAW CAP_SETUID CAP_SETGID

# Memory & OOM hardening
MemoryDenyWriteExecute=true
OOMScoreAdjust=500

[Install]
WantedBy=multi-user.target
