# legacy/debian/postinst
#!/bin/bash
# postinst script for tailscale-receiver

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new-version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <package-being-installed> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

SERVICE_FILE="/usr/lib/systemd/system/tailscale-receive.service"
CONFIG_FILE="/etc/default/tailscale-receive"

case "$1" in
    configure)
        # Configure the service for the current user if possible
        if [ -n "$SUDO_USER" ] && id "$SUDO_USER" >/dev/null 2>&1; then
            TARGET_USER="$SUDO_USER"
        elif [ -n "$USER" ] && [ "$USER" != "root" ] && id "$USER" >/dev/null 2>&1; then
            TARGET_USER="$USER"
        else
            # Try to find a non-root user
            TARGET_USER=$(getent passwd | awk -F: '$3 >= 1000 && $1 != "nobody" {print $1}' | head -1)
            if [ -z "$TARGET_USER" ]; then
                echo "WARNING: Could not determine target user. Please configure TARGET_USER in $CONFIG_FILE"
                TARGET_USER=""
            fi
        fi

        if [ -n "$TARGET_USER" ]; then
            TARGET_DIR="/home/$TARGET_USER/Downloads/tailscale"

            # Update configuration file
            if [ -f "$CONFIG_FILE" ]; then
                sed -i "s|^# TARGET_USER=.*|TARGET_USER=$TARGET_USER|" "$CONFIG_FILE"
                sed -i "s|^# TARGET_DIR=.*|TARGET_DIR=$TARGET_DIR|" "$CONFIG_FILE"
            fi

            # Create target directory
            if [ ! -d "$TARGET_DIR" ]; then
                mkdir -p "$TARGET_DIR"
                chown "$TARGET_USER:$TARGET_USER" "$TARGET_DIR"
                chmod 755 "$TARGET_DIR"
            fi

            # Update systemd service file with user-specific paths
            if [ -f "$SERVICE_FILE" ]; then
                sed -i "s|^ReadWritePaths=.*|ReadWritePaths=/home/$TARGET_USER/Downloads/tailscale|" "$SERVICE_FILE"
            fi

            echo "Configured Tailscale Receiver for user: $TARGET_USER"
            echo "Files will be saved to: $TARGET_DIR"
        fi

        # Reload systemd and enable service
        systemctl daemon-reload || true

        # Only enable/start if Tailscale is available
        if command -v tailscale >/dev/null 2>&1; then
            systemctl enable tailscale-receive.service || true
            echo "Service enabled. Start with: sudo systemctl start tailscale-receive.service"
        else
            echo "Tailscale not found. Please install Tailscale and then start the service:"
            echo "  sudo systemctl start tailscale-receive.service"
        fi

        # Copy Nautilus script to user directory if Nautilus is available
        if command -v nautilus >/dev/null 2>&1 && [ -n "$TARGET_USER" ]; then
            NAUTILUS_DIR="/home/$TARGET_USER/.local/share/nautilus/scripts"
            if [ -f "/etc/skel/.local/share/nautilus/scripts/Send to device using Tailscale" ]; then
                mkdir -p "$NAUTILUS_DIR"
                cp "/etc/skel/.local/share/nautilus/scripts/Send to device using Tailscale" "$NAUTILUS_DIR/"
                chown -R "$TARGET_USER:$TARGET_USER" "$NAUTILUS_DIR"
                chmod +x "$NAUTILUS_DIR/Send to device using Tailscale"
                echo "Nautilus integration installed for user: $TARGET_USER"
            fi
        fi
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

exit 0